# deploy/markdown-studio.service
[Unit]
Description=Студия постов — веб-редактор (server.py)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/markdown-bot
ExecStart=/opt/markdown-bot/venv/bin/python3 app/server.py
Restart=always
RestartSec=5
# наружу смотрит nginx (deploy/nginx-stands.conf), стенд — на внутреннем порту
Environment=PORT=8082
User=markdown
Group=markdown

[Install]
WantedBy=multi-user.target
